# cookbooks/fabmodules/templates/fablocal.service.erb
[Unit]

[Service]
User=fabmodules
Group=fabmodules
WorkingDirectory=/tmp/
ExecStart=/bin/bash /usr/local/bin/fabmodules-local.sh

[Install]
WantedBy=default.target
