# cookbooks/guest_setup/templates/captive_portal.conf.erb
start on desktop-start\nstop on desktop-end\nexec firefox https://wel-clearp1.massey.ac.nz/guest/aruba_guest_wireless.php
